# build/pkg_deb11/supervisord.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
childlogdir=/var/log/supervisor            ; ('AUTO' child log dir, default $TEMP)
nodaemon=true

[program:rsyslog]
command = /usr/sbin/rsyslogd -n -f /etc/rsyslog.conf
numprocs=1
autostart=true
autorestart=true

[program:weewx]
command = /usr/bin/weewxd /etc/weewx/weewx.conf
numprocs=1
autostart=true
autorestart=true

