# supervisord.conf

[supervisord]
nodaemon=true

[program:sshd]
startsecs=0
command=/usr/sbin/sshd

[program:weewx]
startsecs=0
command=/usr/bin/python /home/weewx/bin/weewxd --daemon --pidfile=/var/run/weewx.pid /home/weewx/weewx.conf

[program:nginx]
startsecs=0
command=/usr/sbin/nginx
