# build/setup_centos8/supervisord.conf

[supervisord]
nodaemon=true

[program:rsyslog]
command = /usr/sbin/rsyslogd -n -f /etc/rsyslog.conf
numprocs=1
autostart=true
autorestart=true

[program:weewx]
startsecs=0
command=/usr/bin/python3 /home/weewx/bin/weewxd --daemon --pidfile=/var/run/weewx.pid /home/weewx/weewx.conf

